# services/rtorrent/run
#!/bin/sh

rm -f /config/session/rtorrent.lock
exec script -q -c "rtorrent -n -o import=/etc/rtorrent.rc" /dev/null > /dev/null
2>&1