# services/rtorrent/run
#!/bin/sh

rm -f /config/rtorrent/.session/rtorrent.lock
exec rtorrent -n -o import=/etc/rtorrent.rc
2>&1